# rootfs/etc/services.d/tinyproxy/finish
#!/command/with-contenv bashio
# shellcheck shell=bash

declare exit_code
readonly exit_code_container=$(</run/s6-linux-init-container-results/exitcode)
readonly exit_code_service="${1}"
readonly exit_code_signal="${2}"
readonly service="tinyproxy"

bashio::log.info \
  "Service ${service} exited with code ${exit_code_service}" \
  "(by signal ${exit_code_signal})"
