# rootfs/etc/services.d/tinyproxy/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info 'Starting the tinyproxy daemon...'

exec tinyproxy -d
